# iothub/nginx-mqtt/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
stopsignal=QUIT

[program:sentel-gateway]
command=/sentel/iotgateway
